# data/url-dispatcher-user-watch.conf.in
description "URL Dispatcher"

start on file FILE=$HOME/.config/url-dispatcher/urls/*.url-list

instance $MATCH

exec @pkglibexecdir@/url-dispatcher-update "$MATCH"
