# data/url-dispatcher-update.conf.in
description "URL Dispatcher Directory Watch"

# Work around for: https://bugs.launchpad.net/upstart/+bug/1302117

umask 0022

# end work around

start on (file FILE=~/.config/url-dispatcher/urls/*.url-dispatcher) or (file FILE=@datadir@/url-dispatcher/urls/*.url-dispatcher) or url-dispatcher-update

instance $MATCH
task

pre-start script
	RUNNING=$(initctl status url-dispatcher | grep start/running 2> /dev/null)

	if [ "x${RUNNING}" = "x" ]; then
		echo "DEBUG: Started before URL Dispatcher, let it try to build the DB"
		sleep 5
	fi
end script

exec @pkglibexecdir@/update-directory "$MATCH"
