# data/url-dispatcher-update.conf.in
description "URL Dispatcher Directory Watch"

# Work around for: https://bugs.launchpad.net/upstart/+bug/1302117

umask 0022

# end work around

start on (file FILE=~/.config/url-dispatcher/urls/*.url-dispatcher) or (file FILE=@datadir@/url-dispatcher/urls/*.url-dispatcher) or url-dispatcher-update

instance $MATCH
task

pre-start script
	URL_DISPATCHER_STATE=$(gdbus call --session --dest com.ubuntu.Upstart --object-path /com/ubuntu/Upstart/jobs/url_2ddispatcher/_ --method org.freedesktop.DBus.Properties.Get com.ubuntu.Upstart0_6.Instance state 2> /dev/null)
	if [ -z $URL_DISPATCHER_STATE ] || [ ${URL_DISPATCHER_STATE} != "(<'running'>,)" ] ; then
		echo "Pausing to wait for URL Dispatcher to start"
		sleep 10
	fi
end script

exec @pkglibexecdir@/update-directory "$MATCH"
