# data/url-dispatcher-update.conf.in
description "URL Dispatcher Directory Watch"

start on (file FILE=~/.config/url-dispatcher/urls/*.url-dispatcher) or (file FILE=@datadir@/url-dispatcher/urls/*.url-dispatcher) or url-dispatcher-update

instance $MATCH

exec @pkglibexecdir@/update-directory "$MATCH"
