# data/url-dispatcher-refresh.conf.in
description "Ensure the URL dispatcher database is up-to-date, likely at session init"

start on started url-dispatcher

emits url-dispatcher-update

# Work around for: https://bugs.launchpad.net/upstart/+bug/1302117

umask 0022

# end work around

pre-start script
# If we're starting with url-dispatcher let's let the rest of the
# system have a chance to settle.
	if [ "$UPSTART_EVENTS" = "started" ] ; then
		sleep 60
	fi
end script

script
	initctl emit url-dispatcher-update "MATCH=$HOME/.config/url-dispatcher/urls/"
	initctl emit url-dispatcher-update "MATCH=@datadir@/url-dispatcher/urls/"
end script
