# data/url-dispatcher-system-watch.conf.in
description "URL Dispatcher System Directory Watch"

start on file FILE=@datadir@/url-dispatcher/urls/

exec @pkglibexecdir@/update-directory "@datadir@/url-dispatcher/urls"
