# data/url-dispatcher-failed.conf
description "URL Dispatcher Recovery"

start on stopping url-dispatcher RESULT=failed
task

script
	rm -rf ~/.cache/url-dispatcher/urls-1.db*
end script
