# data/url-dispatcher.conf.in
description "URL Dispatcher"
author "Ted Gould <ted@ubuntu.com>"

start on started dbus
stop on stopping dbus

respawn

emits application-start

exec @pkglibexecdir@/url-dispatcher
