# data/url-dispatcher-user-watch.conf.in
description "URL Dispatcher User Directory Watch"

start on file FILE=~/.config/url-dispatcher/urls/

exec @pkglibexecdir@/update-directory "$MATCH"
