# data/url-dispatcher.conf.in
description "URL Dispatcher"
author "Ted Gould <ted@ubuntu.com>"

start on started dbus
stop on stopping dbus

respawn

# Work around for: https://bugs.launchpad.net/upstart/+bug/1302117

umask 0022

pre-start script
	if [ -d "$HOME/.cache/url-dispatcher" ] && [ ! -x "$HOME/.cache/url-dispatcher" ] ; then
		echo "Fixing cache directory permissions"
		chmod 700 "$HOME/.cache/url-dispatcher"
	fi
end script

# end work around

emits application-start

exec @pkglibexecdir@/url-dispatcher
